# roles/systemd_slice/templates/unit.slice.j2
[Unit]
Description={{ systemd_slice_description }}
Before=slices.target
